# systemd/script.service
[Unit]
Description=PHP bin template

[Service]
Type=oneshot
ExecStart=/usr/bin/php /opt/php-bin-template/bin/script
