# systemd/script.timer
[Unit]
Description=PHP bin template timer

[Timer]
OnCalendar=*-*-* 01:00:00

[Install]
WantedBy=timers.target
